# balancedns-healthcheck.service
[Unit]
Description=BalanceDNS Health Check Service
After=balancedns.service
Wants=balancedns-healthcheck.timer

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/balancedns-healthcheck
ExecStart=/usr/local/bin/balancedns-healthcheck.sh --quiet
Nice=10
IOSchedulingClass=idle
