# balancedns-healthcheck.timer
[Unit]
Description=BalanceDNS Health Check Timer - Runs health check every 15 seconds
Requires=balancedns.service
After=balancedns.service

[Timer]
OnBootSec=10
OnUnitActiveSec=15
AccuracySec=2
Persistent=true
Unit=balancedns-healthcheck.service

[Install]
WantedBy=timers.target
